# lernstick-speech-dispatcher/debian/mimic3.service
[Unit]
Description=Run Mimic 3 web server

[Service]
ExecStart=/usr/bin/mimic3-server

[Install]
WantedBy=default.target
