# lernstick-menu-xdg/etc/xdg/autostart/lock-kde.desktop
[Desktop Entry]
Type=Application
Encoding=UTF-8
Version=1.0
Name=Lock KDE
Name[de_DE]=KDE sperren
Exec=lock-kde
